# salt-call.cron
*/15 * * * * root salt-call --master=salt.finalduty.me state.highstate &>/dev/null
